# WAMwkspc/devel/etc/catkin/profile.d/05.catkin_make_isolated.bash
/home/dhruv/RobotArm/WAMwkspc/devel/.private/catkin/etc/catkin/profile.d/05.catkin_make_isolated.bash