# WAMwkspc/devel/etc/catkin/profile.d/1.ros_version.sh
/home/dhruv/RobotArm/WAMwkspc/devel/.private/ros_environment/etc/catkin/profile.d/1.ros_version.sh