# WAMwkspc/devel/etc/catkin/profile.d/10.rosconsole.sh
/home/dhruv/RobotArm/WAMwkspc/devel/.private/rosconsole/etc/catkin/profile.d/10.rosconsole.sh