# WAMwkspc/devel/etc/catkin/profile.d/10.rosbuild.sh
/home/dhruv/RobotArm/WAMwkspc/devel/.private/rosbuild/etc/catkin/profile.d/10.rosbuild.sh