# WAMwkspc/devel/.private/ros_environment/etc/catkin/profile.d/1.ros_etc_dir.sh
# generated from ros_environment/env-hooks/1.ros_etc_dir.sh.em

# env variable in develspace
export ROS_ETC_DIR="/home/dhruv/RobotArm/WAMwkspc/devel/.private/ros_environment/etc/ros"
